# new-config/.dotfiles/personal/.config/autostart/redshift.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Redshift
Exec=/usr/bin/redshift
Terminal=false
