# old/config/systemd/user/backup.service
[Unit]
Description=Rsync backup

[Service]
ExecStart=%h/.bin/scripts/backup.sh
